# src/StoneAggregator/packaging/root/DEBIAN/prerm
#!/bin/sh
set -e
echo 'Stopping StoneAggregator...'
systemctl stop stoneaggregator
systemctl disable stoneaggregator
